# ledscape.service.in
[Unit]
Description=LEDScape UDP Receiver
ConditionPathExists=/dev/uio0
Wants=network-online.target
After=network-online.target


[Service]
WorkingDirectory=LEDSCAPE_PATH/
ExecStart=LEDSCAPE_PATH/run-ledscape
KillMode=control-group

[Install]
WantedBy=multi-user.target
